# config/sddm
#!/sbin/openrc-run

description="Simple Desktop Display Manager"

command="/usr/bin/sddm"
command_background="yes"
pidfile="/run/sddm.pid"

depend() {
    need localmount
    after bootmisc consolefont devfs hostname keymaps oscclock termencoding udev dbus
}
